# dashboard/opt/dashboard.service
[Unit]
Description=Dashboard web page for 1721's robot.

[Service]
ExecStart=/opt/dashboard/bin/entrypoint.sh

[Install]
WantedBy=multi-user.target
